# boards/mellow-fly-sht-36-v3/98-mellow-fly-sht-36.rules
# Mellow FLY-SHT 36 V3
SUBSYSTEMS=="usb", ATTRS{idProduct}=="614e",  ATTRS{idVendor}=="1d50", ATTRS{serial}=="mellow-fly-sht-36-v3", ACTION=="add", SYMLINK+="mellow-fly-sht-36-v3", RUN+="/home/pi/printer_data/config/RatOS/scripts/klipper-mcu-added.sh"

ACTION=="remove", ENV{DEVLINKS}=="/dev/mellow-fly-sht-36-v3"
